# crontab
*/10 * * * * cd /app/data/public; php bin/worker.php
